# CoSMo-Singapore-ChickenRice.crontab
59 20 * * * ~/DEVEL/GITROOT/DashbordScripts/CoSMo-Singapore-ChickenRice-ITK-runtests.sh > runtests.log

