# CoSMo-Singapore-Durian.crontab
59 20 * * * /Users/CoSMoTest-Mac/DEVEL/GITROOT/DashboardScripts/CoSMo-Singapore-Durian-ITK-runtest.sh > /tmp/log.tests

