# CoSMo-Singapore-ChickenRice.crontab
59 12 * * * ~/DEVEL/GITROOT/DashboardScripts/CoSMo-Singapore-ChickenRice-ITK-runtests.sh > runtests.log

