# SIgN-Karoshi2.crontab
30 8 * * * /Users/alexgouaillard/DEVEL/GITROOT/DashboardScripts/SIgN-Karoshi2-runtests.sh

