# SIgN-Karoshi2.crontab
00 8 * * * $HOME/DEVEL/GITROOT/DashboradScripts/SIgN-Karoshi2-runtests.sh

