# SIgN-Karoshi2.crontab
00 5 * * * /Users/alexgouaillard/DEVEL/GITROOT/DashboardScripts/SIgN-Karoshi2-runtests.sh

