# SIgN-Karoshi2.crontab
00 8 * * * $HOME/DEVEL/GITROOT/DashboardScripts/SIgN-Karoshi2-runtests.sh

